# fcgi.ini
[program:fcgiwrap-cleanup]
command=/bin/rm -f /run/fcgi/fcgi.sock
autostart=true
autorestart=false
startretries=0
priority=1
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:fcgiwrap]
command=/usr/bin/fcgiwrap -f -s unix:/run/fcgi/fcgi.sock
user=nagios
autostart=true
autorestart=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
